# pi-setup/new-skyline-server.service
[Unit]
Description=New Skyline server
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
User=__USER__
WorkingDirectory=__APPDIR__
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=HOST=0.0.0.0
Environment=DATA_SOURCE=api
ExecStart=__PNPM__ -F server start
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
